# templates/razor.erb
#! /bin/bash

### BEGIN INIT INFO
# Provides:          razor
# Required-Start:    $network $mongodb $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Razor 
# Description:       Razor provisioning service
### END INIT INFO

# Author: Stephen Johnson <stephen@puppetlabs.com>

set -e

DAEMON=<%= @directory %>/bin/razor_daemon.rb
test -x $DAEMON || exit 0

. /lib/lsb/init-functions

export PATH="${PATH:+$PATH:}/usr/sbin:/sbin"

case "$1" in
  start)
        log_daemon_msg "Starting razor daemon" "razor"
        $DAEMON start
    ;;
  stop)
        log_daemon_msg "Stopping razor daemon" "razor"
        $DAEMON stop
        log_end_msg $?
    ;;
  restart)
        set +e
        log_daemon_msg "Restarting razor daemon" "razor"
        $DAEMON stop
        sleep 1
        $DAEMON start
    ;;
   status)
        $DAEMON status 
    ;;
  *)
    echo "Usage: /etc/init.d/razor {start|stop|restart|status}"
    exit 1
esac

exit 0
